# files/etc/uci-defaults/99-wifi-opt
uci set wireless.@wifi-device[0].country='CN'
uci set wireless.@wifi-device[0].channel='1'
uci set wireless.@wifi-device[0].htmode='HE20'
uci set wireless.@wifi-device[0].txpower='20'

uci set wireless.@wifi-device[1].country='CN'
uci set wireless.@wifi-device[1].channel='36'
uci set wireless.@wifi-device[1].htmode='HE80'
uci set wireless.@wifi-device[1].txpower='23'

uci commit wireless
